# luci/themes/openmeshbr/root/etc/uci-defaults/luci-theme-openmeshbr
#!/bin/sh
uci batch <<-EOF
	set luci.themes.openmeshbr=/luci-static/openmeshbr
	commit luci
EOF
exit 0
